# templates/profile.d/sbt.sh
#!/bin/bash

# create a symlink to cache dirs
